# packages/tkmm/tkmm_desktop/system.d/tkmm.service
[Unit]
Description=TKMM
Requires=graphical.target

[Service]
ExecStart=/storage/.local/share/tkmm/bin/tkmm
# keep KillMode=process unless there is no good reason to switch to cgroup
KillMode=process
TimeoutStopSec=10
Restart=always
RestartSec=2
StartLimitInterval=0
Type=notify

[Install]
WantedBy=tkmm.target
