# packages/tkmm/tkmm/system.d/tkmm-splash.service
[Unit]
Description=Apply TKMM splash wallpaper
DefaultDependencies=no
Requires=xorg.service
After=xorg.service
Before=graphical.target tkmm-autostart.service tkmm.service
PartOf=tkmm.service

[Service]
Type=oneshot
User=TKMM
Group=TKMM
Environment="DISPLAY=:0"
Environment="HOME=/storage"
ExecStart=/usr/bin/apply-splash

[Install]
WantedBy=graphical.target
WantedBy=tkmm.target
WantedBy=tkmm.service
